# src/.config/autostart/login.sh.desktop
[Desktop Entry]
Exec=/home/tpsapp/.local/bin/login.sh
Icon=application-x-shellscript
Name=login.sh
Type=Application
X-KDE-AutostartScript=true
